# alpine-s6-serf-nginx/rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

s6-setuidgid nginx
fdmove -c 2 1
nginx -c /etc/nginx/nginx.conf
