# alpine-s6-serf-nghttpx-acmesh/rootfs/etc/services.d/nghttpx-access-log/run
#!/usr/bin/execlineb -P

logutil-service -f /var/run/s6/nghttpx-access-log-fifo /var/log/nghttpx-access
