# alpine-s6-serf-nghttpx-acmesh/rootfs/etc/services.d/nghttpx/run
#!/usr/bin/execlineb -P

# nghttpx logs to STDERR. Redirect it to STDOUT so ./log service can track it.
fdmove -c 2 1
nghttpx --conf=/etc/nghttpx/nghttpx.conf
